# systemd/user/theme-monitor.service
[Unit]
Description=Monitor System Theme Changes (Kitty/Neovim)
After=graphical-session.target

[Service]
ExecStart=/usr/bin/python3 %h/.dotfiles/scripts/monitor_theme.py
Restart=on-failure
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=default.target
